# deployment/fusca.service
[Unit]
Description=Fusca Pro Lookup - Auction Data Search
After=network.target

[Service]
User=ubuntu
Group=ubuntu
WorkingDirectory=/var/www/fusca
Environment="PATH=/var/www/fusca/venv/bin"
EnvironmentFile=/etc/fusca-env.conf
ExecStart=/var/www/fusca/venv/bin/gunicorn --workers 1 --bind 127.0.0.1:5001 --timeout 300 app:app
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target

